# templates/unifi-rebuild-keystore.path.j2
[Path]
PathChanged={{ hostvars['localhost']['cert_for_keystore_path'] }}

[Install]
WantedBy=multi-user.target
